# systemd/user/location.timer
[Unit]
Description=Check the current location of your computing

[Timer]
OnActiveSec=30
OnUnitActiveSec=600

[Install]
WantedBy=ready.target
